# target/linux/imx6ull/base-files/etc/uci-defaults/97_dhcp
!/bin/sh
#
# Copyright (C) 2013-2015 OpenWrt.org
#

. /lib/functions/uci-defaults.sh

board=$(board_name)

board_config_update

case "$board" in
lorawan_gateway_ethernet |\
lorawan_gateway_ethernet_bridged)
	uci delete network.lan
	uci delete network.@device[0]
	uci delete network.globals
	uci set network.lan=interface
	uci set network.lan.proto='dhcp'
	uci set network.lan.device='eth0'
	/etc/init.d/network restart
	;;

audio_stream_ethernet)
	uci delete network.@device[0]
	uci set network.@device[0]=device	
	uci set network.@device[0].name='br-lan'
	uci set network.@device[0].type='bridge'
	uci add_list network.@device[0].ports='eth0'
	uci add_list network.@device[0].ports='eth1'
	uci set network.lan.proto='static'
	uci set network.lan.ipaddr='192.168.1.10'
	uci set network.lan.netmask='255.255.255.0'
	;;
esac

board_config_flush

exit 0

