# package/esp32/python-ha-esp32/files/esp32_ble
#!/bin/sh /etc/rc.common
# Copyright (C) 2007 OpenWrt.org

#start after dbus (60)
START=61

boot() {
	hciattach -s 921600 /dev/ttymxc1 any 921600 flow
	hciconfig hci0 up
}
