# target/linux/imx6ull/base-files/etc/hotplug.d/iface/99-wifi-watchdog
#!/bin/sh
if [ "${ACTION}" = "ifdown" -a "${INTERFACE}" = "wwan" ]; then
	/etc/wifi-watchdog.sh "${INTERFACE}"
fi
